# infra/backup/crontab
15 8 * * * PGPASSWORD=$POSTGRES_PASSWORD /usr/local/bin/backup.sh
